# scripts/auto-update.service
[Unit]
Description=TIME Coin Masternode Auto-Update
Documentation=https://github.com/time-coin/time-masternode
# Run after the network is up so git fetch works.
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash /root/time-masternode/scripts/auto-update.sh

# Run as root so update.sh can call systemctl stop/start.
User=root

# Ensure $HOME is set so the script can find ~/time-masternode.
Environment=HOME=/root

StandardOutput=journal
StandardError=journal
SyslogIdentifier=time-auto-update
